# rpcapd@.service
[Unit]
Description=Rpcap Per-Connection Server

[Service]
ExecStart=/usr/bin/rpcapd -n -i
StandardInput=socket
